# autostart/guake.desktop
[Desktop Entry]
Type=Application
Exec=guake
Hidden=false
NoDisplay=false
Name=Guake
Comment=
